# spamassasin_docker/supervisor.conf
[supervisord]
user=root
nodaemon=true


[program:spamassasin]
autostart = true
autorestart = true
command = ./spamd.sh

